# setup/systemd/device-inventory.service
[Unit]
Description=Device Inventory (FastAPI)
After=network.target

[Service]
Type=simple
User=device-inventory
WorkingDirectory=/opt/device-inventory
ExecStart=/opt/device-inventory/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8004 --no-access-log
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
